# packages/atri-main/debian/atri-main/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload || true
systemctl enable atri-main.service || true
