# packages/atri-main/debian/atri-main/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop atri-main.service || true
systemctl disable atri-main.service || true
