# packages/atri-led/lib/systemd/system/atri-led-stop-boot.service
[Unit]
Description=Stop Boot Animation on Login Screen
After=graphical-session-pre.target
Before=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/atrild stop
RemainAfterExit=yes

[Install]
WantedBy=graphical-session-pre.target